# profiles/server/skel/etc/init.d/depmod
#!/sbin/openrc-run

description="Executes depmod command before modprobes in /etc/init.d/modules."

depend() {
	before modules
}

start() {
	depmod
}
